# kiosk.service
[Unit]
Description=Chromium Kiosk
RequiresMountsFor=/run

# Wait for network. Comment the next two lines out if you serve a LOCAL page
# (no internet needed) so the kiosk starts faster / works offline.
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
User=root

# Best-effort clock sync before launch (no-op if offline). Tolerant: a failure
# here never blocks the kiosk from starting.
ExecStartPre=-/bin/sh -c 'datesync 2>/dev/null || /opt/Pi-tools/datesync 2>/dev/null || true'

# Launch arguments come from kiosk.url on the boot partition
# (/boot/firmware/kiosk.url on Raspberry Pi OS, /boot/kiosk.url on x86).
# Run `kiosk --help` to see all options. <mac> is replaced by the drive id.
ExecStart=/bin/sh -c 'F=/boot/firmware/kiosk.url; [ -f "$$F" ] || F=/boot/kiosk.url; KIOSK_ARGS=$$(/usr/bin/grep -v "^#" "$$F" | /usr/bin/tr "\n" " " | /usr/bin/sed -e "s/<mac>/$$(/bin/cat /data/var/drive-id 2>/dev/null)/"); exec /usr/local/bin/kiosk $${KIOSK_ARGS}'
ExecStop=/usr/bin/pkill Xorg
Restart=always

[Install]
WantedBy=multi-user.target
